# crontabs/root
*/10 * * * * /usr/bin/php-cgi /www/lannuestra/php/check.php